# factor_research/scripts/ops/com.astcok.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.astcok.web</string>

  <!-- 常驻 Next.js 前端::3000(生产构建 next start;改前端代码后需 npm run build 重新构建) -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>start</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kiki/astcok/web</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/kiki/astcok/factor_research/logs/web/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kiki/astcok/factor_research/logs/web/launchd.err.log</string>
</dict>
</plist>
