# factor_research/scripts/ops/com.factor.decay_monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.factor.decay_monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/kiki/astcok/factor_research/scripts/ops/decay_monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kiki/astcok/factor_research</string>

    <!-- 每周一 17:00(收盘后,价量已更新;LOOP_ENGINEERING §5.4 周度衰减复测) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>1</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.factor.decay_monitor.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.factor.decay_monitor.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
